# recipes-sota/rvi-sota-client/rvi-sota-client/rvi-sota-client.service
[Unit]
Description=RVI SOTA Client

[Service]
ExecStart=/usr/bin/sota_client -c /var/sota/client.toml

[Install]
WantedBy=multi-user.target
